# debian/loggermcp.service
[Unit]
Description=loggerMCP - MCP server for Ubuntu syslog
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/loggermcp /etc/loggermcp/config.yaml
Restart=on-failure
RestartSec=5
User=loggermcp
Group=loggermcp
WorkingDirectory=/var/lib/loggermcp
StateDirectory=loggermcp
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
ReadOnlyPaths=/var/log

[Install]
WantedBy=multi-user.target
