# debian/postinst
#!/bin/bash
set -euo pipefail

SERVICE_NAME="loggermcp"
CONFIG_DIR="/etc/loggermcp"
CONFIG_FILE="${CONFIG_DIR}/config.yaml"
EXAMPLE_CONFIG_FILE="${CONFIG_DIR}/config.yaml.example"
STATE_DIR="/var/lib/loggermcp"

generate_secret() {
    od -An -N24 -tx1 /dev/urandom | tr -d ' \n'
    echo
}

detect_manifest_hostname() {
    local raw_host manifest_host

    raw_host="$(hostname -f 2>/dev/null || hostname 2>/dev/null || true)"
    manifest_host="$(printf '%s' "$raw_host" | tr '[:upper:]' '[:lower:]' | sed 's/[^a-z0-9.-]/-/g; s/^-*//; s/-*$//')"

    if [ -z "$manifest_host" ]; then
        manifest_host="localhost"
    fi

    printf '%s\n' "$manifest_host"
}

detect_server_ip() {
    local server_ip

    if command -v ip > /dev/null 2>&1; then
        server_ip="$(ip route get 1.1.1.1 2>/dev/null | awk '{for (i = 1; i <= NF; i++) if ($i == "src") {print $(i + 1); exit}}')"
    fi

    if [ -z "${server_ip:-}" ]; then
        server_ip="$(hostname -I 2>/dev/null | awk '{for (i = 1; i <= NF; i++) if ($i ~ /^[0-9]+\./) {print $i; exit}}')"
    fi

    if [ -z "${server_ip:-}" ]; then
        server_ip="127.0.0.1"
    fi

    printf '%s\n' "$server_ip"
}

case "$1" in
    configure)
        # Create service user
        if ! getent passwd "$SERVICE_NAME" > /dev/null 2>&1; then
            adduser --system --group --no-create-home --shell /usr/sbin/nologin "$SERVICE_NAME"
        fi

        # Add to adm group for syslog access
        usermod -aG adm "$SERVICE_NAME" || true

        install -d -m 0750 -o "$SERVICE_NAME" -g "$SERVICE_NAME" "$STATE_DIR"

        # Copy example config if config doesn't exist
        if [ ! -f "$CONFIG_FILE" ]; then
            cp "$EXAMPLE_CONFIG_FILE" "$CONFIG_FILE"

            access_key="$(generate_secret)"
            encryption_key="$(generate_secret)"
            manifest_hostname="$(detect_manifest_hostname)"
            server_ip="$(detect_server_ip)"

            sed -i \
                -e "s|^access_key: .*|access_key: \"${access_key}\"|" \
                -e "s|^encryption_key: .*|encryption_key: \"${encryption_key}\"|" \
                -e "s|^# public_base_url: .*|public_base_url: \"https://${server_ip}:7777\"|" \
                -e "s|^public_base_url: .*|public_base_url: \"https://${server_ip}:7777\"|" \
                -e "s|^manifest_name: .*|manifest_name: \"logger.${manifest_hostname}/mcp\"|" \
                -e "s|^manifest_title: .*|manifest_title: \"loggerMCP (${manifest_hostname})\"|" \
                "$CONFIG_FILE"

            chown "$SERVICE_NAME":"$SERVICE_NAME" "$CONFIG_FILE"
            chmod 600 "$CONFIG_FILE"
            echo "Generated access key, encryption key, manifest name, and public base URL in ${CONFIG_FILE}"
        fi

        if command -v systemctl > /dev/null 2>&1; then
            systemctl daemon-reload
            systemctl enable "${SERVICE_NAME}.service"
        fi
        ;;
esac

exit 0
