# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create service user
        if ! getent passwd loggermcp > /dev/null 2>&1; then
            adduser --system --group --no-create-home --shell /usr/sbin/nologin loggermcp
        fi

        # Add to adm group for syslog access
        usermod -aG adm loggermcp || true

        # Copy example config if config doesn't exist
        if [ ! -f /etc/loggermcp/config.yaml ]; then
            cp /etc/loggermcp/config.yaml.example /etc/loggermcp/config.yaml
            chown loggermcp:loggermcp /etc/loggermcp/config.yaml
            chmod 600 /etc/loggermcp/config.yaml
            echo "WARNING: Edit /etc/loggermcp/config.yaml and set access_key"
        fi

        systemctl daemon-reload
        systemctl enable loggermcp.service
        ;;
esac

exit 0
