# debian/prerm
#!/bin/bash
set -euo pipefail

SERVICE_NAME="loggermcp"

case "$1" in
    remove|purge)
        if command -v systemctl > /dev/null 2>&1; then
            systemctl stop "${SERVICE_NAME}.service" || true
            systemctl disable "${SERVICE_NAME}.service" || true
            systemctl daemon-reload || true
        fi
        ;;
esac

exit 0
